# streaming-dictation/systemd/streaming-dictation.service
[Unit]
Description=Streaming Dictation Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/voicenotebot/streaming-dictation/backend
ExecStart=/usr/bin/node dist/index.js
# Always revive: a stray signal or crash must not leave the service down.
Restart=always
RestartSec=5
# Guard against tight crash loops.
StartLimitIntervalSec=60
StartLimitBurst=5
EnvironmentFile=/root/voicenotebot/streaming-dictation/backend/.env

[Install]
WantedBy=multi-user.target
